# deploy/soulpull-telegram-bot.service
[Unit]
Description=Soulpull Telegram Bot (long-polling)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/soulpull/soulpull-app
EnvironmentFile=/home/soulpull/soulpull-app/.env
ExecStart=/home/soulpull/soulpull-app/venv/bin/python /home/soulpull/soulpull-app/bot/telegram_bot.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target


